# test/overlay-example/overlay-demo/overlay-app.container
[Container]
ContainerName=overlay-app
Image=docker.io/library/nginx:alpine
Pod=overlay-demo.pod
[Install]
WantedBy=default.target
